# support/orm-core-worker.ini
[program:orm-core-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/nginx/OpenResourceManager/Core/artisan queue:work --sleep=5 --tries=15
autostart=true
autorestart=true
user=nginx
numprocs=8
redirect_stderr=true
stdout_logfile=/var/log/orm/core/worker.log